# prev_scripts/update-mywebsite.timer
[Unit]
Description=Timer for My Website Update Service

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
